# cron/crontab
@reboot daemon -- {project}/cron/grabber_downloader.sh
*/10 * * * * php {project}/cron/posts_sched.php >> {project}/logs/cron_posts_sched.log
*/15 * * * * php {project}/cron/vk_to_telegram.php >> {project}/logs/cron_vk_to_telegram.log
*/15 * * * * php {project}/cron/grabber.php >> {project}/logs/cron_grabber.log
0 * * * * php {project}/cron/calc_money.php >> {project}/logs/cron_money.log
*/10 * * * * php {project}/cron/members.php 1 >> {project}/logs/cron_members.log
